# target/init_scripts/systemd-mkl.service
[Unit]
Description=Start the mklauncher announcement program
After=syslog.target network.target

[Service]
Type=forking
ExecStart=/bin/bash /opt/scripts/boot/mkl.sh start
ExecStop=/bin/bash /opt/scripts/boot/mkl.sh stop
User=mdadmn

[Install]
WantedBy=multi-user.target

